# deploy/twelve.service
[Unit]
Description=twelve — workspace walker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=uprootiny
Group=uprootiny
WorkingDirectory=/home/uprootiny/twelve
Environment=TWELVE_PORT=47312
# Bind 0.0.0.0 so the service is reachable both via caddy on :443 (loopback)
# and directly at twelve.solvulator.com:47312 once the GCP firewall opens 47312.
Environment=TWELVE_BIND=0.0.0.0
Environment=TWELVE_ROOT=/home/uprootiny
Environment=TWELVE_HUB2=hub2
ExecStart=/usr/bin/python3 /home/uprootiny/twelve/twelve.py
Restart=on-failure
RestartSec=2

# Hardening — twelve only reads, never writes outside its WorkingDirectory
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/uprootiny/twelve
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
